# crontab
0 10,18 * * * /usr/local/bin/python /app/sync.py
